# corenlp/supervisord.conf
[supervisord]
nodaemon=true

[program:corenlp]
command=java -cp "*" -mx4g edu.stanford.nlp.pipeline.StanfordCoreNLPServer -timeout 15000
directory=/corenlp
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
